# scripts/rc.local
#!/bin/sh -e
echo "Link Start! \n$(date)"> /tmp/rc.local.log
su - tencent -c 'cd ~/RV-AutoAim && screen -dmS "Autoaim" /home/tencent/miniforge3/bin/python3.8 app/aim_process.py production red'
exit 0
